# editions/shared/systemd/specteros-guardian.service
[Unit]
Description=PhantomKernel Guardian Advisor Service
After=phantomkernel-policyd.service
Requires=phantomkernel-policyd.service

[Service]
Type=simple
User=phantomkernel
ExecStart=/usr/bin/phantomkernel-guardian
Restart=on-failure

[Install]
WantedBy=phantomkernel.target

