# editions/debian/packaging/debian/specteros-guardian.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable specteros-guardian.service || true
    systemctl start specteros-guardian.service || true
fi

#DEBHELPER#

exit 0
