# editions/debian/packaging/debian/specteros-policyd.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable specteros-policyd.service || true
    systemctl start specteros-policyd.service || true
fi

#DEBHELPER#

exit 0
