# editions/debian/packaging/debian/specteros-auditd.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable specteros-auditd.service || true
    systemctl start specteros-auditd.service || true
fi

#DEBHELPER#

exit 0
